# 1.3.1/installation/pastell-upstart.conf
description "Pastell Daemon Upstart Script"

start on started mysql
stop on stopping mysql

respawn

exec sudo -u www-data /usr/bin/php /home/eric/adullact/pastell/batch/action-automatique.php